# systemd/mouse-battery-monitor.service
[Unit]
Description=Mouse Battery Monitor
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mouse-battery-monitor.sh
User=root
Environment=HOME=/root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
